# src/etc/services.d/webui/run
#!/usr/bin/with-contenv bash
exec s6-applyuidgid -u ${APP_UID} -g ${APP_GID} \
  /bin/busybox httpd -f -p 6801 -h /webui-aria2-master
