# src/etc/services.d/webui/run
#!/usr/bin/execlineb -P
/bin/busybox httpd -f -p 6801 -h /webui-aria2-master

